# insect-detect-startup.service
[Unit]
Description=Insect Detect Startup Service
Documentation=https://maxsitt.github.io/insect-detect-docs/
Wants=NetworkManager.service
After=network.target NetworkManager.service
Before=multi-user.target

[Service]
Type=oneshot
User=pi
Group=pi
WorkingDirectory=/home/pi/insect-detect
ExecStart=/home/pi/env_insdet/bin/python3 /home/pi/insect-detect/startup.py
Environment=PYTHONPATH=/home/pi/insect-detect
Environment=HOME=/home/pi
RemainAfterExit=yes
Restart=no
TimeoutStartSec=3600
TimeoutStopSec=30
PrivateTmp=true
StandardOutput=journal
StandardError=journal
SyslogIdentifier=insect-detect-startup

[Install]
WantedBy=multi-user.target