# insect-detect-startup.service
[Unit]
Description=Insect Detect Startup Service
Documentation=https://maxsitt.github.io/insect-detect-docs/
Wants=network-online.target
After=network.target network-online.target NetworkManager.service

[Service]
Type=oneshot
# Template file: pi username and /home/pi paths are automatically replaced with the
# actual username and home directory when installing via install.sh or update.sh
User=pi
Group=pi
WorkingDirectory=/home/pi/insect-detect
ExecStart=/home/pi/.local/bin/uv run startup
Environment=HOME=/home/pi
RemainAfterExit=yes
TimeoutStartSec=600
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=insect-detect-startup

[Install]
WantedBy=multi-user.target